# .config/systemd/user/gopls.service
[Unit]
Description=gopls

[Service]
EnvironmentFile=-/home/yuzx/.config/systemd/user/env.conf
ExecStartPre=-env
ExecStart=/home/yuzx/go/bin/gopls serve -listen=127.0.0.1:60098 -logfile=/tmp/gopls.log -debug=127.0.0.1:60099
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=default.target
